# docker/etc/supervisord.d/supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command = /usr/sbin/nginx
user = root
autostart = true

[program:rsyslog]
command = /usr/sbin/rsyslogd -n
user = root
autostart = true

[program:telegraf]
command = /usr/bin/telegraf
user = root
autostart = true

[program:sshd]
command = /usr/sbin/sshd -D -u0
user = root
autostart = true
